# airootfs/etc/skel/.config/autostart/unclutter.desktop
[Desktop Entry]
Type=Application
Name=unclutter
Exec=unclutter &
StartupNotify=false
Terminal=false
